# dns-resolver-cache/unbound-persist-cache.cron
*/5 * * * * root unbound-control dump_cache > /var/cache/unbound.dump~ && mv /var/cache/unbound.dump~ /var/cache/unbound.dump
